# deploy/ec2/daily-jobs.timer
[Unit]
Description=Run daily market data jobs at 06:00 UTC

[Timer]
OnCalendar=*-*-* 06:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
